# Archer_C7_V5/openwrt/build_dir/target-mips_r2_uClibc-0.9.33.2/root-ar71xx/usr/lib/opkg/info/kmod-qca-ssdk-hnat.list
/usr/sbin/qca-hnat
/etc/modules.d/50-qca-ssdk-hnat
/lib/modules/3.3.8/qca-ssdk.ko
/etc/init.d/qca-hnat-start
