# Smart Home > Security Cameras > EC60/image/ALL/rootfs/usr/lib/opkg/info/dnsmasq.list
/lib/upgrade/keep.d/dnsmasq
/etc/dnsmasq.conf
/etc/hotplug.d/iface/25-dnsmasq
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
