# Smart Home > Security Cameras > EC60/image/ALL/rootfs/usr/lib/opkg/info/logModel.list
/usr/sbin/logModel.sh
/usr/bin/syslogd
/etc/init.d/logModelInit
/etc/syslog.conf
/usr/bin/logCtrl
/usr/lib/liblogApi.so
/usr/sbin/logApi
/usr/sbin/logLevel.sh
