# Business > Omada SDN > Routers > ER7206/build_dir/target-mips64-octeon-linux-gnu-osg_pro_un_v1/root-octeon/usr/lib/opkg/info/base-files.list
/etc/wizard_config/config_wan3/service
/etc/wizard_config/config_wan1/ipsec
/etc/wizard_config/config_wan4/dyn3322ddns
/etc/openwrt_version
/etc/wizard_config/config_wan1/ipsec_check_dns_tmp
/etc/wizard_config/config_wan2/time_mngt
/etc/wizard_config/config_wan2/network
/etc/wizard_config/config_wan3/tddp
/etc/wizard_config/config_wan4/ipsec_secrets
/etc/wizard_config/config_wan1/ipsec_pr
/etc/wizard_config/config_wan3/webfilter_global
/etc/wizard_config/config_wan2/ipsec_check_dns
/etc/wizard_config/config_wan2/freeStrategy
/etc/wizard_config/config_wan4/ipsec_check_dns
/etc/wizard_config/config_wan3/url_filter
/etc/resolv.conf
/etc/wizard_config/config_wan2/dos_defense
/etc/wizard_config/config_wan3/logger
/etc/wizard_config/config_wan3/user-secrets.reference
/etc/wizard_config/config_wan2/upnp
/etc/wizard_config/config_wan3/nat
/etc/wizard_config/config_wan3/print_server
/etc/wizard_config/config_wan2/online
/etc/wizard_config/config_wan3/ipgroup
/etc/wizard_config/config_wan2/system_params
/etc/init.d/done
/etc/wizard_config/config_wan1/isp
/var
/etc/wizard_config/config_wan2/pptp-client-global
/etc/wizard_config/config_wan2/l2tp-client
/etc/wizard_config/config_wan1/locale
/etc/nsswitch.conf
/etc/wizard_config/config_wan2/firewall
/etc/wizard_config/config_wan2/policy_route
/etc/wizard_config/config_wan3/user-secrets
/etc/wizard_config/config_wan4/webfilter
/etc/wizard_config/config_wan2/websort
/etc/wizard_config/config_wan4/l2tp-server
/etc/wizard_config/config_wan2/pptp-server
/etc/init.d/system
/etc/preinit
/etc/wizard_config/config_wan2/cmxddns
/etc/wizard_config/config_wan3/switch
/etc/wizard_config/config_wan1/ipstat
/etc/wizard_config/config_wan1/static_route
/etc/wizard_config/config_wan4/user.db
/etc/wizard_config/config_wan3/imb
/etc/wizard_config/config_wan3/splitaccess
/etc/wizard_config/config_wan2/wportal
/lib/upgrade/platform.sh
/etc/wizard_config/config_wan3/reference
/etc/wizard_config/config_wan4/usermngr
/etc/wizard_config/config_wan4/arp_scan_range
/etc/wizard_config/config_wan4/web_security
/etc/shells
/etc/wizard_config/config_wan2/remote_mngt
/etc/wizard_config/config_wan1/webfilter
/etc/wizard_config/config_wan3/accountmgnt
/etc/wizard_config/config_wan2/system
/lib64
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/wizard_config/config_wan2/ipsec_secrets
/etc/wizard_config/config_wan2/luci
/usr/lib64
/etc/wizard_config/config_wan2/dropbear
/etc/wizard_config/config_wan2/ippool
/etc/protocols
/etc/wizard_config/config_wan4/appdist
/etc/uci-defaults/12_network-generate-ula
/etc/wizard_config/config_wan1/phddns
/etc/wizard_config/config_wan4/nat
/etc/wizard_config/config_wan3/pptp-server-global
/etc/wizard_config/config_wan4/load_balance
/etc/wizard_config/config_wan2/isp
/etc/wizard_config/config_wan4/tddp
/etc/wizard_config/config_wan4/flood_defense
/lib/preinit/10_indicate_preinit
/etc/wizard_config/config_wan3/pptp-client
/etc/wizard_config/config_wan2/arp_scan_range
/etc/init.d/smp
/etc/wizard_config/config_wan2/ipsec
/etc/wizard_config/config_wan2/mac_filter
/etc/wizard_config/config_wan4/bwlist_qq
/etc/wizard_config/config_wan4/uhttpd
/etc/wizard_config/config_wan1/line_backup
/etc/wizard_config/config_wan2/arp_defense
/etc/wizard_config/config_wan1/ifstat-mini
/etc/wizard_config/config_wan3/upnp
/etc/wizard_config/config_wan3/dnsproxy
/etc/wizard_config/config_wan3/l2tp-client.reference
/etc/wizard_config/config_wan3/ucitrack
/etc/wizard_config/config_wan4/imb
/etc/wizard_config/config_wan1/print_server
/etc/wizard_config/config_wan4/portal.db
/etc/wizard_config/config_wan2/wireless
/etc/wizard_config/config_wan4/pptp-client-global
/etc/wizard_config/config_wan4/service
/usr/sbin/anete
/etc/wizard_config/config_wan1/bwlist_qq
/etc/wizard_config/config_wan3/access_ctl
/etc/wizard_config/config_wan4/session_limits
/etc/config/system
/etc/wizard_config/config_wan2/dyn3322ddns
/etc/wizard_config/config_wan3/nwadditional
/etc/wizard_config/config_wan4/logger
/etc/wizard_config/config_wan1/user-secrets
/etc/wizard_config/config_wan1/mwan3
/etc/wizard_config/config_wan3/ipsec_check_dns_tmp
/etc/wizard_config/config_wan3/online
/etc/wizard_config/config_wan3/administration
/etc/wizard_config/config_wan3/network
/etc/wizard_config/config_wan2/locale
/etc/wizard_config/config_wan1/usermngr
/usr/sbin/pre_upgrade
/etc/wizard_config/config_wan1/zone
/etc/fstab
/etc/wizard_config/config_wan2/l2tp-server.reference
/etc/sysctl.conf
/etc/wizard_config/config_wan1/portal.db
/sbin/sysupgrade
/etc/wizard_config/config_wan4/ipgroup
/etc/wizard_config/config_wan3/isp
/etc/wizard_config/config_wan2/flood_defense
/lib/preinit/99_10_run_init
/etc/wizard_config/config_wan3/luci
/etc/wizard_config/config_wan4/switch
/etc/wizard_config/config_wan2/session_limits
/etc/wizard_config/config_wan2/ipstat
/etc/wizard_config/config_wan2/web_security
/etc/wizard_config/config_wan2/ipsec_pr
/etc/wizard_config/config_wan4/static_route
/sbin/led.sh
/etc/wizard_config/config_wan4/line_backup
/etc/wizard_config/config_wan2/user-secrets.reference
/lib/preinit/50_indicate_regular_preinit
/etc/wizard_config/config_wan1/administration
/etc/wizard_config/config_wan3/time_mngt
/etc/banner
/etc/wizard_config/config_wan2/l2tp-server
/etc/wizard_config/config_wan3/websort
/etc/wizard_config/config_wan4/ifstat-mini
/etc/cfgsync.d/00_start_sync.sh
/etc/wizard_config/config_wan1/user.db
/etc/wizard_config/config_wan3/cmxddns
/etc/shadow
/etc/wizard_config/config_wan1/url_filter
/etc/wizard_config/config_wan3/firewall
/etc/wizard_config/config_wan3/filter_global
/etc/rc.common
/etc/wizard_config/config_wan1/splitaccess
/etc/uci-defaults/10_migrate-shadow
/etc/wizard_config/config_wan4/upnp
/etc/wizard_config/config_wan2/load_balance
/etc/wizard_config/config_wan3/system
/etc/wizard_config/config_wan3/wportal
/lib/preinit/10_indicate_failsafe
/etc/wizard_config/config_wan3/ippool
/etc/wizard_config/config_wan3/ipsec
/init
/etc/wizard_config/config_wan4/url_filter
/etc/wizard_config/config_wan2/phddns
/etc/wizard_config/config_wan1/accountmgnt
/etc/wizard_config/config_wan1/appdist
/etc/wizard_config/config_wan1/dhcp
/etc/wizard_config/config_wan3/dos_defense
/usr/sbin/setsmp.sh
/etc/wizard_config/config_wan4/isp
/etc/wizard_config/config_wan1/pptp-client
/etc/wizard_config/config_wan3/l2tp-client
/etc/wizard_config/config_wan3/dropbear
/etc/wizard_config/config_wan4/reference
/lib/functions/switch_port.sh
/lib/preinit/80_mount_root
/etc/wizard_config/config_wan1/qos
/etc/wizard_config/config_wan3/freeStrategy
/etc/wizard_config/config_wan2/zone
/etc/wizard_config/config_wan1/nwadditional
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/wizard_config/config_wan2/webfilter
/etc/wizard_config/config_wan4/print_server
/etc/wizard_config/config_wan2/mwan3
/etc/wizard_config/config_wan3/pptp-server
/etc/wizard_config/config_wan1/pptp-client-global
/etc/wizard_config/config_wan4/luci
/etc/wizard_config/config_wan1/filter_global
/etc/init.d/led
/sbin/restorefactory
/etc/wizard_config/config_wan3/policy_route
/etc/wizard_config/config_wan4/splitaccess
/lib/functions/network.sh
/etc/init.d/sysctl
/etc/wizard_config/config_wan4/user-secrets
/etc/wizard_config/config_wan4/online
/etc/wizard_config/config_wan3/system_params
/etc/wizard_config/config_wan4/dnsproxy
/sbin/hotplug-call
/etc/wizard_config/config_wan1/reference
/etc/wizard_config/config_wan1/service
/etc/wizard_config/config_wan3/locale
/etc/wizard_config/config_wan2/webfilter_global
/etc/wizard_config/config_wan2/static_route
/etc/wizard_config/config_wan4/ucitrack
/etc/wizard_config/config_wan3/remote_mngt
/etc/wizard_config/config_wan4/accountmgnt
/usr/sbin/snete
/etc/cfgsync.d/99_end_sync.sh
/etc/wizard_config/config_wan2/l2tp-client.reference
/etc/wizard_config/config_wan4/network
/etc/wizard_config/config_wan3/wireless
/etc/wizard_config/config_wan1/ipsec_check_dns
/etc/wizard_config/config_wan3/ipsec_check_dns
/etc/wizard_config/config_wan3/ipstat
/etc/wizard_config/config_wan3/mac_filter
/etc/rc.button/failsafe
/etc/wizard_config/config_wan1/uhttpd
/etc/wizard_config/config_wan1/ipgroup
/etc/wizard_config/config_wan4/pptp-client
/etc/services
/etc/wizard_config/config_wan1/access_ctl
/etc/wizard_config/config_wan2/usermngr
/etc/wizard_config/config_wan2/qos
/etc/wizard_config/config_wan4/webfilter_global
/etc/wizard_config/config_wan4/ipsec
/etc/wizard_config/config_wan2/bwlist_qq
/etc/wizard_config/config_wan3/arp_defense
/etc/wizard_config/config_wan2/line_backup
/etc/wizard_config/config_wan2/dhcp
/lib/functions/system.sh
/etc/wizard_config/config_wan2/ifstat-mini
/etc/wizard_config/config_wan2/pptp-server-global
/etc/wizard_config/config_wan4/system
/etc/wizard_config/config_wan4/websort
/etc/wizard_config/config_wan3/ipsec_secrets
/etc/uci-defaults/01_network
/etc/wizard_config/config_wan4/cmxddns
/etc/wizard_config/config_wan2/user.db
/lib/upgrade/common.sh
/etc/wizard_config/config_wan4/ippool
/etc/wizard_config/config_wan3/ipsec_pr
/etc/wizard_config/config_wan4/access_ctl
/etc/wizard_config/config_wan3/zone
/etc/wizard_config/config_wan1/logger
/etc/wizard_config/config_wan1/user-secrets.reference
/etc/wizard_config/config_wan1/l2tp-server.reference
/etc/wizard_config/config_wan3/phddns
/etc/wizard_config/config_wan2/portal.db
/etc/wizard_config/config_wan1/system_params
/etc/sysupgrade.conf
/etc/rc.button/power
/etc/wizard_config/config_wan4/wportal
/etc/wizard_config/config_wan1/freeStrategy
/usr/sbin/dnetb
/usr/sbin/dnete
/etc/wizard_config/config_wan4/firewall
/etc/wizard_config/config_wan3/dyn3322ddns
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/etc/wizard_config/config_wan2/print_server
/etc/wizard_config/config_wan2/ipsec_check_dns_tmp
/etc/wizard_config/config_wan1/dnsproxy
/etc/wizard_config/config_wan4/time_mngt
/etc/wizard_config/config_wan3/mwan3
/etc/wizard_config/config_wan3/pptp-client-global
/etc/wizard_config/config_wan1/switch
/etc/wizard_config/config_wan3/arp_scan_range
/etc/wizard_config/config_wan1/ucitrack
/etc/wizard_config/config_wan2/appdist
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/wizard_config/config_wan1/policy_route
/etc/wizard_config/config_wan1/dos_defense
/etc/rc.button/reset
/etc/wizard_config/config_wan4/nwadditional
/etc/wizard_config/config_wan4/l2tp-server.reference
/etc/wizard_config/config_wan2/user-secrets
/lib/preinit/30_failsafe_wait
/etc/wizard_config/config_wan1/l2tp-client
/etc/wizard_config/config_wan3/qos
/etc/wizard_config/config_wan2/url_filter
/sbin/firstboot
/etc/wizard_config/config_wan1/time_mngt
/etc/wizard_config/config_wan1/pptp-server
/etc/wizard_config/config_wan3/l2tp-server
/etc/wizard_config/config_wan1/tddp
/etc/wizard_config/config_wan1/ipsec_secrets
/bin/login.sh
/etc/wizard_config/config_wan4/dropbear
/etc/wizard_config/config_wan1/arp_scan_range
/etc/wizard_config/config_wan4/locale
/etc/wizard_config/config_wan3/web_security
/etc/hotplug.d/net/00-sysctl
/etc/wizard_config/config_wan2/splitaccess
/etc/hosts
/etc/wizard_config/config_wan3/dhcp
/etc/wizard_config/config_wan4/administration
/lib/net_share.sh
/lib/functions.sh
/etc/wizard_config/config_wan3/flood_defense
/etc/wizard_config/config_wan2/service
/etc/wizard_config/config_wan4/ipstat
/etc/wizard_config/config_wan2/uhttpd
/etc/uci-defaults/11_migrate-sysctl
/etc/wizard_config/config_wan4/zone
/etc/wizard_config/config_wan1/remote_mngt
/etc/wizard_config/config_wan2/accountmgnt
/lib/functions/leds.sh
/etc/wizard_config/config_wan4/user-secrets.reference
/etc/wizard_config/config_wan3/webfilter
/etc/wizard_config/config_wan3/load_balance
/etc/wizard_config/config_wan4/pptp-server-global
/etc/wizard_config/config_wan1/network
/etc/wizard_config/config_wan1/firewall
/etc/wizard_config/config_wan3/session_limits
/rom/note
/etc/wizard_config/config_wan4/filter_global
/etc/inittab
/usr/bin/iptables.debug
/etc/openwrt_release
/etc/wizard_config/config_wan4/dos_defense
/usr/sbin/all_ifs
/etc/wizard_config/config_wan1/nat
/etc/wizard_config/config_wan2/administration
/etc/wizard_config/config_wan2/ipgroup
/sbin/wifi
/etc/wizard_config/config_wan2/reference
/etc/wizard_config/config_wan4/wireless
/etc/wizard_config/config_wan4/l2tp-client
/etc/wizard_config/config_wan2/pptp-client
/etc/wizard_config/config_wan4/qos
/etc/wizard_config/config_wan1/arp_defense
/etc/profile
/etc/wizard_config/config_wan2/logger
/etc/wizard_config/config_wan4/phddns
/etc/wizard_config/config_wan1/online
/etc/wizard_config/config_wan1/websort
/etc/wizard_config/config_wan4/mwan3
/etc/wizard_config/config_wan1/l2tp-client.reference
/etc/wizard_config/config_wan4/pptp-server
/etc/wizard_config/config_wan1/cmxddns
/etc/wizard_config/config_wan1/session_limits
/etc/wizard_config/config_wan3/usermngr
/etc/wizard_config/config_wan3/user.db
/etc/init.d/led_early
/etc/wizard_config/config_wan1/mac_filter
/etc/wizard_config/config_wan1/imb
/etc/wizard_config/config_wan4/freeStrategy
/lib/preinit/99_10_failsafe_login
/etc/TZ
/etc/wizard_config/config_wan2/nwadditional
/etc/wizard_config/config_wan1/wportal
/usr/bin/ipxd
/etc/init.d/umount
/etc/wizard_config/config_wan1/flood_defense
/etc/wizard_config/config_wan2/tddp
/lib/functions/service.sh
/etc/wizard_config/config_wan1/dropbear
/usr/bin/ipset.debug
/lib/functions/userconfig.sh
/etc/wizard_config/config_wan2/switch
/etc/wizard_config/config_wan3/bwlist_qq
/etc/wizard_config/config_wan4/ipsec_pr
/lib/upgrade/keep.d/base-files-essential
/usr/bin/ipxr
/etc/wizard_config/config_wan4/mac_filter
/etc/wizard_config/config_wan1/dyn3322ddns
/etc/wizard_config/config_wan4/policy_route
/etc/passwd
/etc/wizard_config/config_wan4/dhcp
/etc/wizard_config/config_wan1/upnp
/etc/wizard_config/config_wan4/remote_mngt
/lib/preinit/02_default_set_state
/usr/sbin/vnetb
/etc/wizard_config/config_wan4/ipsec_check_dns_tmp
/usr/sbin/vnete
/usr/sbin/remote_login_user
/etc/wizard_config/config_wan4/l2tp-client.reference
/etc/mtab
/etc/wizard_config/config_wan1/web_security
/etc/wizard_config/config_wan2/access_ctl
/etc/wizard_config/config_wan3/appdist
/etc/wizard_config/config_wan2/filter_global
/etc/wizard_config/config_wan3/static_route
/etc/init.d/boot
/etc/wizard_config/config_wan2/nat
/lib/functions/octeon.sh
/etc/wizard_config/config_wan1/system
/etc/wizard_config/config_wan3/portal.db
/etc/wizard_config/config_wan2/dnsproxy
/etc/wizard_config/config_wan4/system_params
/etc/wizard_config/config_wan1/ippool
/etc/wizard_config/config_wan2/ucitrack
/lib/functions/preinit.sh
/etc/wizard_config/config_wan3/line_backup
/etc/wizard_config/config_wan4/arp_defense
/bin/ipcalc.sh
/etc/wizard_config/config_wan1/load_balance
/etc/wizard_config/config_wan1/l2tp-server
/etc/wizard_config/config_wan3/ifstat-mini
/etc/wizard_config/config_wan1/wireless
/etc/diag.sh
/etc/wizard_config/config_wan1/webfilter_global
/etc/wizard_config/config_wan2/imb
/etc/wizard_config/config_wan3/uhttpd
/etc/wizard_config/config_wan3/l2tp-server.reference
/etc/wizard_config/config_wan1/pptp-server-global
/etc/wizard_config/config_wan1/luci
