# Business > Omada SDN > Routers > ER7206/build_dir/target-mips64-octeon-linux-gnu-osg_pro_un_v1/root-octeon/usr/lib/opkg/info/lldpd.list
/etc/config/lldpd
/lib/lldpd/lldp_get_wan_device.sh
/usr/lib/liblldpctl.so.4.8.0
/usr/sbin/lldpd
/usr/sbin/lldpcli
/usr/lib/liblldpctl.so
/etc/init.d/lldpd
/usr/lib/liblldpctl.so.4
/usr/sbin/lldpctl
