# Business > Omada SDN > Routers > ER7206/build_dir/target-mips64-octeon-linux-gnu-osg_pro_un_v1/root-octeon/usr/lib/opkg/info/nginx.list
/www/luci-static/resources/authserver/mobile_wifi.html
/etc/nginx/nginx.conf
/etc/nginx/fastcgi_params
/www/luci-static/resources/authserver/pc_wifi.html
/etc/nginx/koi-win
/www/luci-static/resources/authserver/wechat.html
/www/luci-static/resources/authserver/pcauth.js
/usr/sbin/nginx
/etc/nginx/win-utf
/etc/nginx/koi-utf
/etc/cfgsync.d/65_nginx_sync.sh
/www/luci-static/resources/authserver/wifi-bgd.jpg
/www/luci-static/resources/authserver/style-simple-follow.css
/etc/init.d/nginx
/www/luci-static/resources/authserver/style-pcdemo.css
/etc/nginx/mime.types
/lib/upgrade/keep.d/nginx
/www/luci-static/resources/authserver/web_login.html
/usr/lib/lua/ngx_authserver/ngx_wdas.lua
/www/luci-static/resources/authserver/nobody.jpg
/etc/hotplug.d/bridgescripts/99-nginx.sh
/usr/lib/lua/ngx_authserver/ngx_sqlApi.lua
/www/luci-static/resources/authserver/background.jpg
/usr/lib/lua/ngx_authserver/ngx_init.lua
/usr/lib/lua/ngx_authserver/authlistCheck.lua
/www/luci-static/resources/authserver/wifi-logo.png
