# Business > Omada SDN > Routers > ER7206/build_dir/target-mips64-octeon-linux-gnu-osg_pro_un_v1/root-octeon/usr/lib/opkg/info/openvpn-openssl.list
/lib/openvpn/openvpn-common.sh
/lib/openvpn/openvpn-server-up.sh
/etc/config/openvpn
/lib/openvpn/openvpn-client-up.sh
/usr/sbin/checkpsw.sh
/etc/init.d/openvpn
/lib/openvpn/openvpn-client-down.sh
/usr/sbin/openvpn
/lib/openvpn/openvpn-instance.sh
/lib/openvpn/openvpn-server-down.sh
/lib/upgrade/keep.d/openvpn
