# Business > Omada SDN > Routers > ER7206/build_dir/target-mips64-octeon-linux-gnu-osg_pro_un_v1/root-octeon/usr/lib/opkg/info/snmpd.list
/etc/init.d/snmpd
/etc/default/snmpd
/usr/sbin/snmpd
